# project/main.cron
30 2 * * * /home/de10/rdkv/project/main.py &> /home/de10/rdkv/project.log 
